# packaging/rpm/sources/profile.d/volta.csh
setenv VOLTA_HOME /usr/lib/volta
setenv PATH "$VOLTA_HOME/bin:$PATH"
