# packaging/rpm/sources/profile.d/volta.sh
export VOLTA_HOME=/usr/lib/volta
export PATH=$VOLTA_HOME/bin:$PATH
